# supervisord.conf
[supervisord]
nodaemon=true

[program:php5-fpm]
command=/bin/bash -c "service php5-fpm start"

[program:nginx]
command=/bin/bash -c "service nginx start"